# docker/php/supervisor/conf.d/default.ini
[program:queue-worker-default]
directory=/app
command=/usr/local/bin/php artisan queue:listen
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/var/log/supervisor-queue-worker-default.log
autostart=false
autorestart=true
user=root
stopsignal=TERM
numprocs=%(ENV_DEFAULT_WORKERS_COUNT)s
